# conf/supervisor/ghost.conf
[program:ghost]
command = play /opt/ansible/ghost.yml runas ghost npm start
directory = /ghost
priority = 500
autostart = true
startsecs = 10
environment = HOME="/ghost"
